# backend/configuration/supervisor.conf
[supervisord]
nodaemon=true

[program:rethinkdb]
command=/usr/bin/rethinkdb --bind all
startsecs=60